# udev/09.ttyAMA0.rules
KERNEL=="ttyAMA0", ACTION=="add", RUN+="/usr/bin/stty -F /dev/ttyAMA0 -echo"
